# resources/supervisor.conf
[program:teashka.localzet.com]
user = root
command = php master restart
directory = /var/www/teashka
numprocs = 1
autorestart = true
autostart = true